# etc/xdg/autostart/touchegg.desktop
[Desktop Entry]
Exec=touchegg
Icon=system-run
Path=
Terminal=false
Type=Application
